# data/zeitgeist.service.in
[Unit]
Description=Zeitgeist activity log service

[Service]
Type=dbus
BusName=org.gnome.zeitgeist.Engine
ExecStart=@prefix@/bin/zeitgeist-daemon
